# speedtest/cron/crontab
* * * * * bash /date-script.sh >> /var/log/myjob.log 2>&1
