# speedtest/cron/crontab
* * * * * python3 /app/speedtest.py >> /var/log/speedtest.log 2>&1
